# bin/crontab
SHELL="/bin/bash"
HOME="/home/plugin"
PATH="/home/plugin/bin:/home/plugin/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

@reboot /bin/bash $HOME/bin/start.sh >>$HOME/log/boot.log 2>&1

# Keepalive
# mm hh         dom mon dow command
*/5 *            *   *   *  /bin/bash $HOME/bin/keepalive.sh >>$HOME/log/keepalive.log 2>&1

# Restarts
# mm hh         dom mon dow command
02 */8           *   *   *  /bin/bash $HOME/bin/restart.sh >>$HOME/log/restart.log 2>&1

# Cleanup
# mm hh         dom mon dow command
07 07            *   *   *  /bin/bash $HOME/bin/cleanup_logs.sh >>$HOME/log/cleanup_logs.log 2>&1

# Backup
# mm hh         dom mon dow command
#31 07            *   *   *  /bin/bash $HOME/bin/backup.sh >>$HOME/log/backup.log 2>&1
